# tailscale/rootfs/etc/services.d/caddy/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Tailscale
# Runs the Caddy daemon
# ==============================================================================
declare ha_service_protocol
declare domain

if bashio::config.false "enable_proxy"; then
  s6-svc -O .
  bashio::exit.ok
fi

bashio::log.info "Starting Caddy..."
# Wait for socket to be available
while ! bashio::fs.socket_exists "/var/run/tailscale/tailscaled.sock";
do
  sleep 1
done

bashio::log.debug "Checking if SSL is used..."
if bashio::var.true "$(bashio::core.ssl)" ; then
    ha_service_protocol="https"
else
    ha_service_protocol="http"
fi
bashio::log.debug "ha_service_protocol: ${ha_service_protocol}"

if bashio::var.is_empty "${ha_service_protocol}" ; then
    bashio::exit.nok "Error checking if SSL is enabled"
fi

while true;
do
  if tailscale status --json --peers=false --self=false \
    | jq --exit-status '.BackendState == "Running"' &> /dev/null;
  then
    IFS=","
    domain=$(tailscale status --self=true --peers=false --json | jq -rc ".CertDomains[]")
    bashio::log.info "Your Homeassistant instance is available via Tailscale."
    bashio::log.info "Your Domain is ${domain}"
    caddy reverse-proxy --from=${domain} --to=${ha_service_protocol}://homeassistant:$(bashio::core.port)
  fi
  # Well... wait a bit more
  sleep 2
done