# tailscale/rootfs/etc/s6-overlay/s6-rc.d/tailscale/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Tailscale
# Runs after the machine has been logged in into the Tailscale network and
# ==============================================================================
declare -a interfaces
declare -a addresses
declare -a routes
declare -a options
declare ipinfo
declare tags

# Find interfaces and matching addresses from which we can extract routes to be advertised
interfaces+=($(bashio::network.interfaces))
for interface in "${interfaces[@]}"; do
    addresses+=($(bashio::network.ipv4_address "${interface}"))
    addresses+=($(bashio::network.ipv6_address "${interface}"))
done
routes=()

for address in "${addresses[@]}"; do
  if bashio::var.has_value "${address}"; then
    # Skip local link addresses
    if [[ "${address:0:6}" == "fe80::" ]] || [[ "${address:0:8}" == "169.254." ]];
    then
      continue
    fi

    # Skip if forwarding for the address family is disabled
    if [[ "${address}" =~ .*:.* ]];
    then
      [[ $(</proc/sys/net/ipv6/conf/all/forwarding) -eq 0 ]] && continue
    else
      [[ $(</proc/sys/net/ipv4/ip_forward) -eq 0 ]] && continue
    fi

    ipinfo="$(/usr/bin/ipcalc --json "${address}")"
    routes+=($(bashio::jq "${ipinfo}" '.NETWORK + "/" + .PREFIX'))
  fi
done

# Get configured tags
if bashio::config.exists "tags" ; then
  tags=$(bashio::config "tags | join(\",\")" "")
else
  tags=""
fi

# Wait for socket to be available
while ! bashio::fs.socket_exists "/var/run/tailscale/tailscaled.sock";
do
  sleep 1
done

# Wait for the network to be available and logged in
while true;
do
  if tailscale status --json --peers=false --self=false \
    | jq --exit-status '.BackendState == "Running" or .BackendState == "NeedsLogin"' &> /dev/null;
  then
    IFS=","

    options+=(--hostname "$(bashio::info.hostname)")
    options+=(--accept-routes)
    options+=(--advertise-tags="${tags}")

    if bashio::config.true "enable_ssh" ; then
      options+=(--ssh=true)
    else
      options+=(--ssh=false)
    fi

    if bashio::config.true "enable_exit_node" ; then
      options+=(--advertise-exit-node=true)
    else
      options+=(--advertise-exit-node=false)
    fi

    if bashio::config.true "enable_subnets" ; then
      options+=(--advertise-routes="${routes[*]}")
    else
      options+=(--advertise-routes="")
    fi

    bashio::log.info 'Starting Tailscale VPN...'
    tailscale up "${options[@]}"
    bashio::exit.ok
  fi
  # Well... wait a bit more
  sleep 2
done