# tailscale/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Tailscale
# Runs the Caddy daemon
# ==============================================================================
declare ha_service_protocol
declare domain

# Wait for socket to be available
while ! bashio::fs.socket_exists "/var/run/tailscale/tailscaled.sock";
do
  sleep 1
done

bashio::log.debug "Checking if SSL is used..."
if bashio::var.true "$(bashio::core.ssl)" ; then
    ha_service_protocol="https"
else
    ha_service_protocol="http"
fi
bashio::log.debug "ha_service_protocol: ${ha_service_protocol}"

if bashio::var.is_empty "${ha_service_protocol}" ; then
    bashio::exit.nok "Error checking if SSL is enabled"
fi

while true;
do
  if tailscale status --json --peers=false --self=false \
    | jq --exit-status '.BackendState == "Running"' &> /dev/null;
  then
    IFS=","
    # Check if Tailnet domain is available
    if tailscale status --self=true --peers=false --json | jq -rce ".CertDomains[]" &> /dev/null ; then
      domain=$(tailscale status --self=true --peers=false --json | jq -rc ".CertDomains[]")
      bashio::log.info "Your Homeassistant instance is available via Tailscale."
      bashio::log.info "Your Domain is https://${domain}"
      # start Caddy proxy for HA instance
      bashio::log.info "Starting Caddy..."
      caddy reverse-proxy --from=${domain} --to=${ha_service_protocol}://homeassistant:$(bashio::core.port)
    else
      bashio::log.warning "Can't get your Tailscale domain, is HTTPS support enabled?"
      bashio::log.warning "Enable HTTPS for your Tailnet to use the caddy proxy"
      bashio::log.warning "https://tailscale.com/kb/1153/enabling-https/"

      # bring service down and don't start it again
      s6-svc -O .
      bashio::exit.ok
    fi
  fi
  # Wait for backend (Tailscaled) to be available and logged in (Tailscale)
  sleep 2
done